# quadlets-meta.target
[Unit]
Description=Meta quadlet target
Wants=quadlets-hostwatch.target

[Install]
WantedBy=default.target