# quadlets-meta.target
[Unit]
Description=Meta quadlet target
Wants=quadlets-hostwatch.target
Wants=quadlets-distrobox.target

[Install]
WantedBy=default.target