# services/langtool/langtool.target
[Unit]
Description=LanguageTool quadlet target

[Install]
WantedBy=default.target