# hostwatch/quadlets-hostwatch.target
[Unit]
Description=hostwatch quadlet target
Upholds=quadlets-hostwatch.path

[Install]
WantedBy=default.target