# services/syncthing/syncthing.target
[Unit]
Description=Syncthing quadlet target

[Install]
WantedBy=default.target