# toolboxes/ephemeral/bluefin-dx-cli.target
[Unit]
Description=bluefin-dx-cli quadlet target

[Install]
WantedBy=default.target