# toolboxes/ephemeral/fedora-distrobox-quadlet.target
[Unit]
Description=fedora-distrobox-quadlet quadlet target

[Install]
WantedBy=default.target