# toolboxes/ephemeral/ubuntu-distrobox-quadlet.target
[Unit]
Description=ubuntu-distrobox-quadlet quadlet target

[Install]
WantedBy=default.target