# toolboxes/quadlets-distrobox.target
[Unit]
Description=Distrobox quadlet target
PartOf=quadlets-meta.target

[Install]
WantedBy=default.target