# toolboxes/quadlets-distrobox.target
[Unit]
Description=Distrobox quadlet target

[Install]
WantedBy=default.target